# transport/coap/target/control/deb/postrm
#!/bin/sh

update-rc.d -f tb-coap-transport remove
