# transport/http/target/debian/postrm
#!/bin/sh -e



update-rc.d -f tb-http-transport remove


