# transport/mqtt/target/control/deb/postinst
#!/bin/sh

chown -R tb-mqtt-transport: /var/log/tb-mqtt-transport
chown -R tb-mqtt-transport: /usr/share/tb-mqtt-transport
update-rc.d tb-mqtt-transport defaults

