# transport/mqtt/target/debian/postrm
#!/bin/sh -e



update-rc.d -f tb-mqtt-transport remove


