# transport/http/target/control/deb/postinst
#!/bin/sh

chown -R tb-http-transport: /var/log/tb-http-transport
chown -R tb-http-transport: /usr/share/tb-http-transport
update-rc.d tb-http-transport defaults

