# transport/http/target/control/rpm/postinst
#!/bin/sh

chown -R tb-http-transport: /var/log/tb-http-transport
chown -R tb-http-transport: /usr/share/tb-http-transport

if [ $1 -eq 1 ] ; then
        # Initial installation
        systemctl --no-reload enable tb-http-transport.service >/dev/null 2>&1 || :
fi
